# scripts/honey-heist.service
[Unit]
Description=Honey Heist Game Server
After=network.target

[Service]
Type=simple
User=corey
WorkingDirectory=/home/corey/honey-heist
ExecStart=/usr/bin/npx serve --listen tcp://0.0.0.0:8000 .
Restart=always
RestartSec=3
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
